# launchd/com.williamhsiao.codex-app-server.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.williamhsiao.codex-app-server</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>__PROJECT_ROOT__/scripts/codex-app-server.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_ROOT__</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>CODEX_APP_SERVER_LISTEN_URL</key>
      <string>ws://127.0.0.1:8765</string>
    </dict>

    <key>RunAtLoad</key>
    <true />

    <key>KeepAlive</key>
    <true />

    <key>StandardOutPath</key>
    <string>/tmp/codex-app-server.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/codex-app-server.error.log</string>
  </dict>
</plist>
